# supervisord-storm-supervisor.ini
[program:storm-supervisor]
command=storm supervisor
directory=/home/storm
autorestart=true
user=storm
startsecs=10
startretries=999
log_stdout=true
log_stderr=true
logfile=/var/log/storm/supervisor.out
logfile_maxbytes=20MB
logfile_backups=10
